# crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
LOGFILE=/app/logs/cron.log

30 7 * * * python /app/cleaner.py >> $LOGFILE 2>&1
*/5 * * * * python /app/export.py >> /app/logs/cron.log 2>&1
00 20 * * * python /app/import.py >> $LOGFILE 2>&1
